# scripts/backup.list
/etc/NetworkManager/NetworkManager.conf
/etc/apt/
/etc/exports
/etc/forked-daapd.conf
/etc/mdm/Init/Default
/etc/minidlna.conf
/etc/network/interfaces
/etc/ppp/chap-secrets
/etc/ppp/pptpd-options
/etc/pptpd.conf
/etc/rc.local
/etc/samba/smb.conf
/etc/ssh/sshd_config
/etc/sysctl.conf
/etc/vsftpd.conf
/etc/yum.repos.d
/etc/yum/
/etc/init.d/
/etc/sysconfig/
/etc/fstab
/etc/sudoers
/etc/sudoers.d/
/home/rraheja/Desktop/
/home/rraheja/Documents/
/home/rraheja/Projects/
/home/rraheja/.ssh/
/home/rraheja/.config/
/home/rraheja/.conky/
/home/rraheja/.lua/
/home/rraheja/.gnome2/
/home/rraheja/.kde/
/home/rraheja/.linuxmint/
/home/rraheja/.nevernote/
/home/rraheja/.vnc/
/home/rraheja/.VirtualBox/
/home/oracle/Desktop/
/home/oracle/Projects/
/home/oracle/.ssh/
/home/oracle/.config/
/home/oracle/.gnome2/
/home/oracle/.kde/
/home/oracle/.vnc/
/home/twins/Desktop/
/home/twins/Documents/
/home/twins/Projects/
